# .config/autostart/alttab.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=alttab
Comment=alttab
Exec=alttab
OnlyShowIn=bspwm