# .config/autostart/clipmenud.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=clipmenud
Comment=clipmenud
Exec=clipmenud
OnlyShowIn=bspwm