# .config/autostart/config-screen.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=config-screen
Comment=config-screen
Exec=bash -c "sh $HOME/bin/config-screen.sh"
OnlyShowIn=bspwm
