# .config/autostart/conky.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=conky
Comment=conky
Exec=conky
OnlyShowIn=bspwm