# .config/autostart/dunst.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=dunst
Comment=dunst
Exec=dunst
OnlyShowIn=bspwm