# .config/autostart/glava.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=glava
Comment=glava
Exec=bash -c "~/bin/run_glava.sh"
OnlyShowIn=bspwm