# .config/autostart/picom.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=picom
Comment=picom
Exec=bash -c "picom --config ~/.config/picom.conf"
OnlyShowIn=bspwm
