# .config/autostart/polybar.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=polybar
Comment=polybar
Exec=bash -c "~/bin/run_polybar.sh"
OnlyShowIn=bspwm