# .config/autostart/sxhkd.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=sxhkd
Comment=sxhkd
Exec=sxhkd
OnlyShowIn=bspwm